# httpdate
#!/sbin/openrc-run
#
#/***************************************************************************
# *   Copyright (C) 2014,2016 Daniel Mueller (deso@posteo.net)              *
# *                                                                         *
# *   This program is free software: you can redistribute it and/or modify  *
# *   it under the terms of the GNU General Public License as published by  *
# *   the Free Software Foundation, either version 3 of the License, or     *
# *   (at your option) any later version.                                   *
# *                                                                         *
# *   This program is distributed in the hope that it will be useful,       *
# *   but WITHOUT ANY WARRANTY; without even the implied warranty of        *
# *   MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the         *
# *   GNU General Public License for more details.                          *
# *                                                                         *
# *   You should have received a copy of the GNU General Public License     *
# *   along with this program.  If not, see <http://www.gnu.org/licenses/>. *
# ***************************************************************************/
#
# $Header: Sync current time over HTTP (works even behind proxy)$

depend() {
  after dns dhcpcd net ntpdate
}

start() {
  if [ -x /etc/init.d/ntpdate ]
  then
    /etc/init.d/ntpdate status > /dev/null
    # if ntpdate started successfull then just do nothing
    # ntpdate status will exit with 0 if started, not 0 if not
    if [ $? -eq 0 ]
    then
      einfo "ntpdate already running, avoiding more imprecise sync"
      return 0
    fi
  fi

  ebegin "Syncing time"
  httpdate="$(wget -S -O /dev/null google.com 2>&1 | sed -n -e 's/  *Date: *//p' -eT -eq)"
  if [ $? -a ! -z "$httpdate" ]
  then
    date --set="$httpdate" > /dev/null
    eend ${?}
  else
    eerror "Failed to retrieve current time"
    eend 1
  fi
}
